# root/etc/cron.daily/maid
#!/usr/bin/env bash
source /home/coaxial/.rvm/environments/ruby-2.6.5@maid
sudo -Hu coaxial bash -li -c "maid clean --force --silent --rules ~/.maid/rules.rb"
